# conf/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx

[program:php]
command=php-fpm7.0 -F

[include]
files=/etc/supervisor/conf.d/*.conf
