# conf/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
childlogdir=/var/log/supervisor
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx

[program:php-fpm]
command=/usr/sbin/php5-fpm -c /etc/php5/fpm -F